# ezdock/98-ezdock.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="ezdock.pl"
ACTION=="remove", SUBSYSTEM=="i2c", RUN+="ezdock.pl"
